# cloak/files/etc/uci-defaults/70-uhttpd
#!/bin/sh
#
# Copyright 2014 Bright Things Unlimited
#
uci batch <<-EOF
	delete uhttpd.main.listen_http
	add_list uhttpd.main.listen_http=10000
	commit uhttpd
EOF
exit 0
